# _run_/quadlets/object-detection-custom-stream-viewer.container
# ~/.config/containers/systemd/object-detection-custom-stream-viewer.container

[Unit]
Description=Stream Viewer — WebSocket server for browser UI
After=object-detection-custom-camera-capture.service
Requires=object-detection-custom-camera-capture.service

[Container]
Image=quay.io/luisarizmendi/object-detection-custom-stream-viewer:latest
ContainerName=object-detection-custom-stream-viewer
PullPolicy=missing

Volume=/dev/shm:/dev/shm
Network=host
SecurityLabelDisable=true

Environment=SHM_FRAME_NAME=camera_frame
Environment=ZMQ_SUB_HOST=localhost
Environment=ZMQ_SUB_PORT=5555
Environment=ZMQ_TOPIC=detections
Environment=HTTP_PORT=8080
Environment=FRAME_JPEG_QUALITY=75
Environment=MAX_STREAM_FPS=35
Environment=LOG_LEVEL=INFO

PublishPort=8080:8080

[Service]
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
