# _run_/quadlets/object-detection-custom-camera-capture.container
# ~/.config/containers/systemd/object-detection-custom-camera-capture.container
# or /etc/containers/systemd/object-detection-custom-camera-capture.container (system-wide)
#
# Auto-detects USB cameras via v4l2-ctl (format/resolution/fps).
# Falls back to looping video files from VID_DIR if no camera is found.
#
# Enable:
#   systemctl --user daemon-reload
#   systemctl --user start object-detection-custom-camera-capture

[Unit]
Description=Camera Capture — USB cam (auto-detect) or video fallback → shared memory
After=network.target

[Container]
Image=quay.io/luisarizmendi/object-detection-custom-camera-capture:latest
ContainerName=object-detection-custom-camera-capture

# 39 is the video group id. If id is different in your system change it here.
UserNS=keep-id:uid=1000,gid=39
PodmanArgs=--group-add=keep-groups

Volume=/dev/shm:/dev/shm
# Uncomment to enable video file fallback:
Volume=/var/home/detector/videos:/videos:ro

# Pass all video devices so auto-detection works.
# Restrict to a specific device if preferred.
AddDevice=/dev/video0
# AddDevice=/dev/video1

Network=host
SecurityLabelDisable=true

# Environment
Environment=CAM_TARGET_FPS=30
Environment=DEVICE_PROBE_TIMEOUT=15
Environment=VID_DIR=/videos
Environment=SHM_FRAME_NAME=camera_frame
Environment=VIEWER_JPEG_QUALITY=80
Environment=LOG_LEVEL=INFO

# Overrides (uncomment to use):
# Environment=CAMERA_DEVICE=/dev/video0
# Environment=CAM_FRAMERATE=30
# Environment=CAM_RESOLUTION=1280x720
# Environment=CAM_FORMAT=MJPG

[Service]
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
